# systemd-docker-cleanup-service/docker-cleanup.timer
[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitInactiveSec=10s

[Install]
WantedBy=timers.target